# a2a-zen.service
[Unit]
Description=Zen A2A Server
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/a2a-server
ExecStart=/usr/bin/python3 /home/ec2-user/a2a-server/server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
